# rootfs/etc/services.d/openresty/run
#!/bin/sh

while [ ! -f /usr/local/openresty/nginx/lock ]; do sleep 1; done

/usr/local/openresty/bin/openresty -g "daemon off;" 2>&1
